# deploy/syncra-api.service
[Unit]
Description=Syncra API
After=network.target redis-server.service
Wants=redis-server.service

[Service]
Type=simple
WorkingDirectory=/var/www/Syncra/backend
ExecStart=/usr/share/dotnet/dotnet Syncra.Api.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=syncra-api
User=www-data
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://localhost:5000

[Install]
WantedBy=multi-user.target
